# internal/graphrag/taxonomy/techniques/mitre_enterprise.yaml
# MITRE ATT&CK Enterprise Techniques
# This file contains key techniques from the MITRE ATT&CK Enterprise matrix.
# Techniques are organized by tactic and include platform information.
# Full taxonomy available at: https://attack.mitre.org

techniques:
  # Initial Access
  - technique_id: "T1190"
    name: "Exploit Public-Facing Application"
    taxonomy: "mitre"
    category: "initial_access"
    tactic: "Initial Access"
    description: "Adversaries may attempt to exploit weaknesses in an Internet-facing host or system to initially access a network"
    platforms: ["Windows", "Linux", "macOS", "Network"]
    examples:
      - "Exploiting SQL injection vulnerabilities"
      - "Exploiting remote code execution vulnerabilities in web applications"
    detection_strategies:
      - "Monitor application logs for exploitation attempts"
      - "Use Web Application Firewalls (WAF)"
      - "Implement intrusion detection systems"
    url: "https://attack.mitre.org/techniques/T1190/"

  - technique_id: "T1566"
    name: "Phishing"
    taxonomy: "mitre"
    category: "initial_access"
    tactic: "Initial Access"
    description: "Adversaries may send phishing messages to gain access to victim systems"
    platforms: ["Windows", "macOS", "Linux", "Office 365", "SaaS", "Google Workspace"]
    examples:
      - "Spearphishing attachment with malicious payload"
      - "Spearphishing link to credential harvesting page"
    detection_strategies:
      - "Monitor for suspicious email attachments"
      - "Analyze email headers and links"
      - "User awareness training"
    url: "https://attack.mitre.org/techniques/T1566/"

  # Execution
  - technique_id: "T1059"
    name: "Command and Scripting Interpreter"
    taxonomy: "mitre"
    category: "execution"
    tactic: "Execution"
    description: "Adversaries may abuse command and script interpreters to execute commands, scripts, or binaries"
    platforms: ["Windows", "macOS", "Linux"]
    examples:
      - "PowerShell execution"
      - "Bash/sh script execution"
      - "Python interpreter abuse"
    detection_strategies:
      - "Monitor process execution logs"
      - "Detect unusual command-line arguments"
      - "Analyze script contents"
    url: "https://attack.mitre.org/techniques/T1059/"

  - technique_id: "T1203"
    name: "Exploitation for Client Execution"
    taxonomy: "mitre"
    category: "execution"
    tactic: "Execution"
    description: "Adversaries may exploit software vulnerabilities in client applications to execute code"
    platforms: ["Windows", "macOS", "Linux"]
    examples:
      - "Exploiting browser vulnerabilities"
      - "Exploiting PDF reader vulnerabilities"
    detection_strategies:
      - "Monitor for unusual application behavior"
      - "Keep software updated and patched"
    url: "https://attack.mitre.org/techniques/T1203/"

  # Persistence
  - technique_id: "T1078"
    name: "Valid Accounts"
    taxonomy: "mitre"
    category: "persistence"
    tactic: "Persistence"
    description: "Adversaries may obtain and abuse credentials of existing accounts to persist in the environment"
    platforms: ["Windows", "macOS", "Linux", "SaaS", "IaaS", "Network", "Containers"]
    examples:
      - "Using stolen credentials"
      - "Creating backdoor accounts"
    detection_strategies:
      - "Monitor account creation and login patterns"
      - "Implement multi-factor authentication"
      - "Audit privileged account usage"
    url: "https://attack.mitre.org/techniques/T1078/"

  - technique_id: "T1053"
    name: "Scheduled Task/Job"
    taxonomy: "mitre"
    category: "persistence"
    tactic: "Persistence"
    description: "Adversaries may abuse task scheduling functionality to facilitate initial or recurring execution of malicious code"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Windows scheduled tasks"
      - "Cron jobs on Linux"
    detection_strategies:
      - "Monitor scheduled task creation"
      - "Audit task scheduler logs"
    url: "https://attack.mitre.org/techniques/T1053/"

  # Privilege Escalation
  - technique_id: "T1068"
    name: "Exploitation for Privilege Escalation"
    taxonomy: "mitre"
    category: "privilege_escalation"
    tactic: "Privilege Escalation"
    description: "Adversaries may exploit software vulnerabilities to elevate privileges"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Kernel exploits"
      - "Local privilege escalation vulnerabilities"
    detection_strategies:
      - "Monitor for unexpected privilege elevation"
      - "Keep systems patched"
    url: "https://attack.mitre.org/techniques/T1068/"

  - technique_id: "T1548"
    name: "Abuse Elevation Control Mechanism"
    taxonomy: "mitre"
    category: "privilege_escalation"
    tactic: "Privilege Escalation"
    description: "Adversaries may circumvent mechanisms designed to control elevate privileges to gain higher-level permissions"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Bypass User Account Control"
      - "Sudo abuse"
    detection_strategies:
      - "Monitor elevation events"
      - "Audit sudo usage"
    url: "https://attack.mitre.org/techniques/T1548/"

  # Defense Evasion
  - technique_id: "T1027"
    name: "Obfuscated Files or Information"
    taxonomy: "mitre"
    category: "defense_evasion"
    tactic: "Defense Evasion"
    description: "Adversaries may attempt to make an executable or file difficult to discover or analyze"
    platforms: ["Windows", "macOS", "Linux"]
    examples:
      - "Code obfuscation"
      - "Encrypted payloads"
      - "Steganography"
    detection_strategies:
      - "Analyze file entropy"
      - "Behavioral analysis"
    url: "https://attack.mitre.org/techniques/T1027/"

  - technique_id: "T1070"
    name: "Indicator Removal"
    taxonomy: "mitre"
    category: "defense_evasion"
    tactic: "Defense Evasion"
    description: "Adversaries may delete or modify artifacts generated on a system to remove evidence of their presence"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Clear Windows Event Logs"
      - "Delete bash history"
    detection_strategies:
      - "Monitor log deletion events"
      - "Implement centralized logging"
    url: "https://attack.mitre.org/techniques/T1070/"

  # Credential Access
  - technique_id: "T1110"
    name: "Brute Force"
    taxonomy: "mitre"
    category: "credential_access"
    tactic: "Credential Access"
    description: "Adversaries may use brute force techniques to gain access to accounts"
    platforms: ["Windows", "macOS", "Linux", "Office 365", "SaaS", "IaaS", "Network"]
    examples:
      - "Password guessing"
      - "Credential stuffing"
    detection_strategies:
      - "Monitor failed login attempts"
      - "Implement account lockout policies"
      - "Rate limiting"
    url: "https://attack.mitre.org/techniques/T1110/"

  - technique_id: "T1003"
    name: "OS Credential Dumping"
    taxonomy: "mitre"
    category: "credential_access"
    tactic: "Credential Access"
    description: "Adversaries may attempt to dump credentials to obtain account login and credential material"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "LSASS memory dumping"
      - "/etc/shadow extraction"
    detection_strategies:
      - "Monitor for LSASS access"
      - "Detect credential dumping tools"
    url: "https://attack.mitre.org/techniques/T1003/"

  # Discovery
  - technique_id: "T1087"
    name: "Account Discovery"
    taxonomy: "mitre"
    category: "discovery"
    tactic: "Discovery"
    description: "Adversaries may attempt to get a listing of valid accounts"
    platforms: ["Windows", "macOS", "Linux"]
    examples:
      - "net user command"
      - "LDAP queries"
    detection_strategies:
      - "Monitor account enumeration commands"
    url: "https://attack.mitre.org/techniques/T1087/"

  - technique_id: "T1046"
    name: "Network Service Discovery"
    taxonomy: "mitre"
    category: "discovery"
    tactic: "Discovery"
    description: "Adversaries may attempt to get a listing of services running on remote hosts"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Port scanning"
      - "Service enumeration"
    detection_strategies:
      - "Monitor for port scanning activity"
      - "Detect network reconnaissance tools"
    url: "https://attack.mitre.org/techniques/T1046/"

  # Lateral Movement
  - technique_id: "T1021"
    name: "Remote Services"
    taxonomy: "mitre"
    category: "lateral_movement"
    tactic: "Lateral Movement"
    description: "Adversaries may use valid accounts to log into a service that accepts remote connections"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "RDP usage"
      - "SSH access"
    detection_strategies:
      - "Monitor remote login attempts"
      - "Analyze connection logs"
    url: "https://attack.mitre.org/techniques/T1021/"

  - technique_id: "T1570"
    name: "Lateral Tool Transfer"
    taxonomy: "mitre"
    category: "lateral_movement"
    tactic: "Lateral Movement"
    description: "Adversaries may transfer tools or other files between systems in a compromised environment"
    platforms: ["Windows", "macOS", "Linux"]
    examples:
      - "File transfers via SMB"
      - "SCP/SFTP usage"
    detection_strategies:
      - "Monitor file transfers"
      - "Analyze network traffic"
    url: "https://attack.mitre.org/techniques/T1570/"

  # Collection
  - technique_id: "T1005"
    name: "Data from Local System"
    taxonomy: "mitre"
    category: "collection"
    tactic: "Collection"
    description: "Adversaries may search local system sources to find files of interest"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "File system searches"
      - "Database queries"
    detection_strategies:
      - "Monitor file access patterns"
    url: "https://attack.mitre.org/techniques/T1005/"

  # Command and Control
  - technique_id: "T1071"
    name: "Application Layer Protocol"
    taxonomy: "mitre"
    category: "command_and_control"
    tactic: "Command and Control"
    description: "Adversaries may communicate using application layer protocols to avoid detection"
    platforms: ["Windows", "macOS", "Linux"]
    examples:
      - "HTTP/HTTPS C2 channels"
      - "DNS tunneling"
    detection_strategies:
      - "Analyze network traffic patterns"
      - "Monitor for unusual protocol usage"
    url: "https://attack.mitre.org/techniques/T1071/"

  - technique_id: "T1573"
    name: "Encrypted Channel"
    taxonomy: "mitre"
    category: "command_and_control"
    tactic: "Command and Control"
    description: "Adversaries may employ a known encryption algorithm to conceal command and control traffic"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "TLS encrypted C2"
      - "Custom encryption"
    detection_strategies:
      - "Analyze encrypted traffic patterns"
      - "Monitor certificate usage"
    url: "https://attack.mitre.org/techniques/T1573/"

  # Exfiltration
  - technique_id: "T1041"
    name: "Exfiltration Over C2 Channel"
    taxonomy: "mitre"
    category: "exfiltration"
    tactic: "Exfiltration"
    description: "Adversaries may steal data by exfiltrating it over an existing command and control channel"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Data exfiltration via HTTP POST"
      - "DNS exfiltration"
    detection_strategies:
      - "Monitor outbound data transfers"
      - "Analyze C2 traffic"
    url: "https://attack.mitre.org/techniques/T1041/"

  # Impact
  - technique_id: "T1486"
    name: "Data Encrypted for Impact"
    taxonomy: "mitre"
    category: "impact"
    tactic: "Impact"
    description: "Adversaries may encrypt data on target systems to interrupt availability"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Ransomware encryption"
    detection_strategies:
      - "Monitor for mass file encryption"
      - "Detect ransomware behavior"
    url: "https://attack.mitre.org/techniques/T1486/"

  - technique_id: "T1499"
    name: "Endpoint Denial of Service"
    taxonomy: "mitre"
    category: "impact"
    tactic: "Impact"
    description: "Adversaries may perform Endpoint Denial of Service attacks to degrade or block availability"
    platforms: ["Windows", "Linux", "macOS"]
    examples:
      - "Resource exhaustion"
      - "Application crashes"
    detection_strategies:
      - "Monitor resource utilization"
      - "Detect DoS patterns"
    url: "https://attack.mitre.org/techniques/T1499/"
